# 推荐开源代码借鉴/Master—四川大学/Inital/scurm.service
[Unit]
Description=SCURM Auto-Load Script

[Service]
Type=simple
TimeoutStartSec=5min
ExecStart=  /root/scurm/scurm.run
ExecReload=/root/scurm/scurm.run
ExecStop=  /root/scurm/scurm.ext
Restart=always

[Install]
WantedBy=multi-user.target
